# debian/pop-shop.postinst
systemctl --global enable io.elementary.appcenter.timer