# data/io.elementary.appcenter.timer
[Unit]
Description=Check for package updates daily

[Timer]
OnCalendar=daily
OnActiveSec=10

[Install]
WantedBy=timers.target
