# data/io.elementary.appcenter-flatpak-repair-system.service
[Unit]
Description=Flatpak Repair
After=network.target network-online.target systemd-networkd.service NetworkManager.service connman.service

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak repair --system

[Install]
WantedBy=multi-user.target
